# debian/alff-agent.postinst
#!/bin/sh -e
#
# alff-agent.postinst
#

if [ -x "/etc/init.d/alff-agent" ]; then
	# start alff-agent as soon as possible after the network has been initialized.
	update-rc.d alff-agent defaults >/dev/null
fi

# create rules dir
. /etc/alff/alff-defaults.conf
[ -d ${RULES_DIR} ] || mkdir -p ${RULES_DIR}
