# debian/alff-agent.postrm
#!/bin/sh -e
#
# alff-agent.postrm
#

if [ "$1" = purge ]; then
	update-rc.d alff-agent remove
fi
